# deploy/fluxbot-dashboard.service
[Unit]
Description=FluxBot Dashboard
After=network-online.target
Wants=network-online.target
# If Postgres runs on this same machine, uncomment the next line so the
# dashboard waits for it on boot. Leave it out if Postgres is remote.
# Requires=postgresql.service

[Service]
Type=simple
User=fluxbot
Group=fluxbot
WorkingDirectory=/opt/fluxbot
ExecStart=/opt/fluxbot/venv/bin/python /opt/fluxbot/run_dashboard.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fluxbot-dashboard

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
